# enterprise/gateway/tenant-router.service
[Unit]
Description=OpenClaw Tenant Router
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu
Environment=HOME=/home/ubuntu
EnvironmentFile=-/etc/openclaw/env
ExecStart=/usr/bin/python3 /home/ubuntu/tenant_router.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
